# AX6-IPQ/files/etc/init.d/ax6-irq-affinity
#!/bin/sh /etc/rc.common
# AX6 IRQ 亲和性 — 把 ath11k / ethernet / NSS 钉到不同核
# START=94: 在 smp_affinity(93) 之后运行, 确保我们的分配为最终生效值
# shellcheck disable=SC2034  # START/STOP read by /etc/rc.common framework
START=94
STOP=10

set_aff() {
    [ -e "/proc/irq/$1/smp_affinity" ] && echo "$2" > "/proc/irq/$1/smp_affinity" 2>/dev/null
}

apply() {
    # 4 核 IPQ8074(避免 NSS UBI 与 wifi 抢同核):
    #   CPU0 (0x1)  系统 / kthreads
    #   CPU1 (0x2)  eth/dp + wifi1 (2.4G)
    #   CPU2 (0x4)  wifi0 (5G,大流量)
    #   CPU3 (0x8)  NSS UBI doorbell / 固件事件(独占核)
    # NSS 内核 thread 自带绑核,不受这里影响
    while IFS=: read -r irq rest; do
        irq=$(echo "$irq" | tr -d ' ')
        case "$rest" in
            # 5G WiFi → CPU2 (kernel 6.18: ppdu-end-interrupts-mac1, 旧: ath11k_hw0)
            *ppdu-end*mac1*|*wifi0*|*ath11k*hw0*) set_aff "$irq" 4 ;;
            # 2.4G WiFi → CPU1 (kernel 6.18: ppdu-end-interrupts-mac3, 旧: ath11k_hw1)
            *ppdu-end*mac3*|*wifi1*|*ath11k*hw1*) set_aff "$irq" 2 ;;
            # Ethernet + NSS DP → CPU1
            *eth*|*edma*|*dp[0-9]*|*nss-dp*) set_aff "$irq" 2 ;;
            # NSS firmware → CPU3 独占
            *nss*) set_aff "$irq" 8 ;;
        esac
    done < /proc/interrupts
}

start() {
    sleep 5
    apply
}

reload() { apply; }
restart() { apply; }
