# AX6-IPQ/files/etc/uci-defaults/95-ax6-nss-vlan-guard
#!/bin/sh
# AX6 NSS VLAN Guard — 移除 bridge_vlan_filtering(与 NSS WiFi 不兼容)
# 来源: https://github.com/qosmio/openwrt-ipq/blob/main-nss/nss-setup/example/README.md
#
# 生效场景:首次启动 / sysupgrade 后,清理可能被带过来的旧配置
#
# NSS 用户**必须**用 8021q 子接口(lan1.20)实现 VLAN,
# 不能用 DSA bridge_vlan_filtering / `config bridge-vlan` / 端口 'lan1:u*' 标签

CHANGED=0

# 清除所有 device 的 vlan_filtering=1
for d in $(uci -q show network 2>/dev/null | awk -F'[.=]' '/=device$/{print $2}'); do
    cur=$(uci -q get network.${d}.vlan_filtering)
    if [ "$cur" = "1" ]; then
        uci -q delete network.${d}.vlan_filtering
        CHANGED=1
        logger -t nss-vlan-guard "Disabled vlan_filtering on network.$d (incompatible with NSS WiFi)"
    fi
done

# 删除任何 bridge-vlan 段(DSA 语法,与 NSS 不兼容)
while uci -q show network 2>/dev/null | grep -q '=bridge-vlan'; do
    section=$(uci -q show network 2>/dev/null | awk -F'[.=]' '/=bridge-vlan/{print $2; exit}')
    [ -z "$section" ] && break
    uci -q delete network.${section}
    CHANGED=1
    logger -t nss-vlan-guard "Removed bridge-vlan section network.$section (incompatible with NSS WiFi)"
done

if [ "$CHANGED" = "1" ]; then
    uci commit network
    logger -t nss-vlan-guard "VLAN config sanitized for NSS compatibility"
fi

exit 0
