# AX6-IPQ/files/etc/uci-defaults/97-ax6-cn-defaults
#!/bin/sh
# AX6 国内友好默认值: NTP + DNS (uci-defaults — 仅首次启动)
# sysupgrade 保留配置时跳过用户自定义项 (逐条检查, 不破坏混合配置)

# ===========================================================
# NTP: 阿里 + 国家授时中心 + 腾讯 (仅当全部为 OpenWrt 默认时替换)
# ===========================================================
apply_ntp=1
if uci -q get system.ntp.server >/dev/null 2>&1; then
    for s in $(uci -q get system.ntp.server 2>/dev/null); do
        case "$s" in
            *openwrt.pool.ntp.org|*openwrt.org) ;;
            ntp.aliyun.com|ntp.tencent.com|ntp.ntsc.ac.cn|cn.pool.ntp.org) ;;
            *) apply_ntp=0; break ;;  # 用户自定义了非默认服务器
        esac
    done
fi
if [ "$apply_ntp" = "1" ]; then
    uci -q delete system.ntp.server
    uci -q add_list system.ntp.server='ntp.aliyun.com'
    uci -q add_list system.ntp.server='ntp.tencent.com'
    uci -q add_list system.ntp.server='ntp.ntsc.ac.cn'
    uci -q add_list system.ntp.server='cn.pool.ntp.org'
    uci -q set system.ntp.enabled='1'
    uci -q set system.ntp.enable_server='0'
    uci -q set system.@system[0].timezone='CST-8'
    uci -q set system.@system[0].zonename='Asia/Shanghai'
    uci -q commit system
fi

# ===========================================================
# DNS 上游: 阿里 + 腾讯 + DNSPod (仅当全部为已知默认时替换)
# ===========================================================
apply_dns=1
if uci -q get dhcp.@dnsmasq[0].server >/dev/null 2>&1; then
    for s in $(uci -q get dhcp.@dnsmasq[0].server 2>/dev/null); do
        case "$s" in
            223.5.5.5|119.29.29.29|114.114.114.114) ;;
            *) apply_dns=0; break ;;  # 用户自定义了非默认 DNS
        esac
    done
fi
if [ "$apply_dns" = "1" ]; then
    uci -q delete dhcp.@dnsmasq[0].server 2>/dev/null
    uci -q add_list dhcp.@dnsmasq[0].server='223.5.5.5'
    uci -q add_list dhcp.@dnsmasq[0].server='119.29.29.29'
    uci -q add_list dhcp.@dnsmasq[0].server='114.114.114.114'
    uci -q set dhcp.@dnsmasq[0].noresolv='1'
    uci -q set dhcp.@dnsmasq[0].rebind_protection='1'
    uci -q set dhcp.@dnsmasq[0].rebind_localhost='1'
    uci -q set dhcp.@dnsmasq[0].localise_queries='1'
    uci -q set dhcp.@dnsmasq[0].domainneeded='1'
    uci -q set dhcp.@dnsmasq[0].cachesize='2000'
    uci -q commit dhcp
fi

exit 0

