# AX6-IPQ/files/etc/hotplug.d/ieee80211/99-ax6-reapply-irq
#!/bin/sh
# wifi reload 后重新绑 IRQ (延迟等待 IRQ 注册, nohup 防止 SIGUP)
[ "$ACTION" = "add" ] && nohup sh -c 'sleep 2; /etc/init.d/ax6-irq-affinity reload' >/dev/null 2>&1 &

