# AX6-IPQ/files/etc/hotplug.d/ieee80211/99-ax6-reapply-irq
#!/bin/sh
# wifi reload 后重新绑 IRQ (nohup 防止 SIGUP 终止后台进程)
[ "$ACTION" = "add" ] && nohup /etc/init.d/ax6-irq-affinity reload >/dev/null 2>&1 &
